# services/deployment/gestionproductos.socket
[Unit]
Description=Socket

[Socket]
ListenStream=/run/gestionproductos.sock

[Install]
WantedBy=sockets.target
